# installer/roles/receptor_install/files/receptor_controller.service
[Unit]
Description=Receptor Relay Service

[Service]
ExecStart=/usr/local/bin/receptor controller

[Install]
WantedBy=multi-user.target