# installer/roles/receptor_install/files/receptor_node.service
[Unit]
Description=Receptor Relay Service

[Service]
ExecStart=/usr/local/bin/receptor node

[Install]
WantedBy=multi-user.target